# templates/etc/systemd/system/keepup.service.epp
[Unit]
Description=Send keepup package versions

[Service]
Type=oneshot
ExecStart=/opt/keepup/run.sh
