# templates/etc/systemd/system/keepup.timer.epp
<%- |
  String $calendar,
| -%>
[Unit]
Description=Run keepup periodically

[Timer]
OnCalendar=<%= $calendar %>
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
